# luci-app-koti/root/etc/uci-defaults/80_koti
#!/bin/sh

touch /etc/config/koti
uci set koti.dns='dns'
uci set koti.dns.enable='1'
uci add_list koti.dns.server='9.9.9.9'
uci commit koti

return 0
